# docs/examples/devdb-hub-watch.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.devdb.hub-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/devdb</string>
    <string>hub-sync</string>
    <string>--watch</string>
    <string>--interval</string>
    <string>60</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/path/to/devdb-or-target-repo</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/devdb-hub-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/devdb-hub-sync.err</string>
</dict>
</plist>
